# scripts/systemd/tv-off.service
[Unit]
Description=Turn TV Off via HDMI CEC
After=network.target

[Service]
Type=oneshot
User=pi
ExecStart=/var/www/kiosk/scripts/tv-control.sh off
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
