# scripts/systemd/unclutter.service
[Unit]
Description=Hide mouse cursor with unclutter
After=default.target
Before=kiosk.service

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/unclutter -idle 0.1 -root
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
